# os-image/ubuntu/24.04-picoclaw/picoclaw.service
[Unit]
Description=PicoClaw AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/picoclaw gateway
Restart=on-failure
RestartSec=5
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
